# packages/clustermgr/clustermgr-3.1.4.post7.tar.gz/clustermgr-3.1.4.post7/clustermgr/templates/twemproxy/twemproxy.service
[Unit]
Description=nutcracker - High performance proxy server for memcached/redis
After=network.target

[Service]
Type=simple
ExecStartPre=/usr/sbin/nutcracker -t -c /etc/nutcracker/nutcracker.yml
ExecStart=/usr/sbin/nutcracker -c /etc/nutcracker/nutcracker.yml -o /var/log/nutcracker/nutcracker.log
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target