# packages/x-mroy-1052/x-mroy-1052-0.1.7.tar.gz/x-mroy-1052-0.1.7/x-auth.conf
[program:x-auth]
command=/usr/local/bin/x-telserver
stdout_logfile= /var/log/x-auth.log 
stderr_logfile= /var/log/x-auth.err.log
