# packages/x-mroy-1052/x-mroy-1052-0.1.7.tar.gz/x-mroy-1052-0.1.7/x-neid.conf
[program:x-neid]
command=/usr/local/bin/x-neid -p 60009
stdout_logfile= /var/log/x-neid.log
stderr_logfile= /var/log/x-neid.err.log
