# packages/x-mroy-1052/x-mroy-1052-0.1.7.tar.gz/x-mroy-1052-0.1.7/x-node-test.conf
[program:x-node-test]
command=/usr/local/bin/x-ea-test -s start
stdout_logfile= /var/log/x-sstest.log 
stderr_logfile= /var/log/x-sstest.err.log
