# packages/ldh-middleware/ldh_middleware-0.0.3.tar.gz/ldh_middleware-0.0.3/config_sample/supervisord/purist_middleware_monitor.conf
; stored as /etc/supervisor/conf.d/purist_middleware_monitor.conf
; Copyright 2017-2018 Purism SPC
; SPDX-License-Identifier: AGPL-3.0-or-later

; Adapted from extra/supervisord/celerybeat.conf in Celery
; Copyright 2012-2014 GoPivotal, Inc.
; Copyright 2009-2012, 2015-2016 Ask Solem and contributors
; https://github.com/celery/celery/blob/master/extra/supervisord/celerybeat.conf
; SPDX-License-Identifier: BSD-3-Clause

[program:purist_middleware_monitor]
command=/opt/purist/middleware_virtualenv/bin/celery worker -B --app middleware --scheduler django_celery_beat.schedulers:DatabaseScheduler --loglevel=INFO
directory=/opt/purist/middleware
user=www-data
numprocs=1
stdout_logfile=/var/log/purist/middleware/beat.log
stderr_logfile=/var/log/purist/middleware/beat.log
autostart=true
autorestart=true
startsecs=10
environment=VIRTUAL_ENV=/opt/purist/middleware_virtualenv

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
